# cookbooks/bootstrap/templates/hab-sup.service.erb
[Unit]
Description=Habitat Supervisor

[Service]
ExecStart=/bin/hab sup run <%= @hab_sup_opts %>

[Install]
WantedBy=default.target